# systemd/dcom-proxy.service
[Unit]
Description=DCOM Proxy Box — Web Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/dcom-proxy
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
