# systemd/ddns-update.service
[Unit]
Description=DDNS Updater for DCOM Proxy Box
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/dcom-proxy/scripts/ddns-update.sh
StandardOutput=journal
StandardError=journal
